# config/mindspace-beta.service
[Unit]
Description=Mindspace Beta Viz Server
After=network.target

[Service]
Type=simple
User=developer
WorkingDirectory=/home/developer/workspaces/github/PichlerThomas/xpollination-mcp-server
ExecStart=/home/developer/.nvm/versions/node/v22.22.0/bin/node viz/server.js 4201
Environment=NODE_ENV=staging
Environment=DATABASE_PATH=/home/developer/workspaces/github/PichlerThomas/xpollination-mcp-server/data/xpollination-beta.db
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
