# systems/hetzner-cx22/configs/systemd/mindspace-test.service
[Unit]
Description=XPollination Mindspace Test Server (port 4200, VPN-only, develop branch)
After=network.target

[Service]
Type=simple
User=developer
WorkingDirectory=/home/developer/workspaces/github/PichlerThomas/xpollination-mcp-server-test
EnvironmentFile=/home/developer/workspaces/github/PichlerThomas/xpollination-mcp-server-test/.env.test
ExecStart=/home/developer/.nvm/versions/node/v22.22.0/bin/node viz/server.js 4200
Restart=on-failure
RestartSec=5
Environment=NODE_ENV=test
Environment=VIZ_BIND=10.33.33.1

[Install]
WantedBy=multi-user.target
